# doc/syslogd.service
[Unit]
Description=Syslog Director

[Service]
TimeoutStartSec=0
Restart=always
ExecStart=/opt/syslogd/bin/syslogd --port 514 --no-stdout --syslog=udp://localhost:1514

[Install]
WantedBy=default.target
