# doc/syslogd.service
[Unit]
Description=Syslog Director

[Service]
TimeoutStartSec=0
Restart=always
ExecStart=/opt/syslogd/bin/syslogd --port 514 --no-ansi

[Install]
WantedBy=default.target
